# Docker/s6/services.d/crond/run
#!/usr/bin/execlineb -P
foreground { echo "starting crond" }
/usr/sbin/crond -f
